# oodt-filemgr/scripts/supervisord.conf
[supervisord]
nodaemon=true

#[program:solr]
#directory=%(ENV_SOLR_INSTALL_DIR)s/bin
#command=sh -c "./solr start -p 8983 -s $SOLR_HOME -d $SOLR_INSTALL_DIR/server -Dsolr.autoSoftCommit.maxTime=5000 -Dsolr.autoCommit.maxTime=10000 -a '-Dsolr.data.dir=%(ENV_SOLR_DATA_DIR)s'"
#stdout_logfile=/dev/stdout
#stdout_logfile_maxbytes=0

[program:oodt_filemgr]
directory=%(ENV_OODT_HOME)s/cas-filemgr/bin
command=sh -c "./filemgr start"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
