# oodt-fmprod/supervisord.conf
[supervisord]
nodaemon=true

[program:tomcat]
directory=%(ENV_CATALINA_HOME)s/bin
command=sh -c "./catalina.sh start"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
