# oodt-wmgr-resmgr/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:oodt_workflow]
directory=%(ENV_OODT_HOME)s/cas-workflow/bin
command=sh -c "./wmgr start"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:oodt_resource]
directory=%(ENV_OODT_HOME)s/cas-resource/bin
command=sh -c "./resmgr start && tail -f ../logs/cas_resource.log"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:oodt_stub]
directory=%(ENV_OODT_HOME)s/cas-resource/bin
command=sh -c "%(ENV_JAVA_HOME)s/bin/java -Djava.ext.dirs=../lib org.apache.oodt.cas.resource.system.extern.XmlRpcBatchStub --portNum 2001"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
